# varnish/supervisord.conf
[supervisord]
nodaemon=true

[program:varnishd]
command=/usr/sbin/varnishd -a 0.0.0.0:6081 -f /etc/varnish/default.vcl
autorestart=true

