# website/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:varnishd]
command=/usr/sbin/varnishd -F -a 0.0.0.0:6081 -f /etc/varnish/default.vcl -s malloc,100M
autorestart=true

